# systemd/udcn-quic.service
[Unit]
Description=µDCN QUIC Transport Router
After=network.target
Documentation=https://github.com/yourusername/rust-udcn

[Service]
Type=simple
User=udcn
Group=udcn
ExecStartPre=/bin/mkdir -p /var/log/udcn
ExecStart=/usr/local/bin/udcn quic router --config /etc/udcn/quic-config.toml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=udcn-quic

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/log/udcn /etc/udcn
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
